# etc/crontabs/root
* * * * * lua /etc/plzwatchdog_cron.lua > /dev/null >/dev/null 2>&1
* * * * * lua /etc/bleloop_cron.lua > /dev/null >/dev/null 2>&1